# supervisord.conf
[supervisord]
nodaemon=true
user=root

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:app1]
command=uvicorn whisperflow.fast_server:app --host 0.0.0.0 --port 8000
directory=/usr/src/app
autorestart=true
stdout_logfile=/var/log/supervisor/app1.log
stderr_logfile=/var/log/supervisor/app1_err.log